# deploy/rcmdd.service
[Unit]
Description=rcmd relay server
Documentation=https://github.com/obay/rcmd
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rcmd
Group=rcmd
# Edit /etc/rcmd/rcmdd.toml before starting; see rcmdd.toml.example.
ExecStart=/usr/local/bin/rcmdd serve --config /etc/rcmd/rcmdd.toml
Restart=always
RestartSec=2
# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/rcmd
# autocert needs to bind :80 and :443
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
